# extra_conf/files/docker-gc.service
[Unit]
Description=Run Docker-GC

[Service]
User=root
Type=oneshot
ExecStart=/usr/local/bin/docker-gc.sh