# extra_conf/files/wpscan.timer
[Unit]
Description=Run wp-scan weekly

[Timer]
OnCalendar=Sat *-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
