# extra_conf/files/wpscan.service
[Unit]
Description=Run WP-Scan

[Service]
User=elatov
Type=oneshot
ExecStart=/usr/local/bin/wp-scan.sh
