# extra_conf/files/wpscan.service
[Unit]
Description=Run WP-Scan

[Service]
User=elatov
Type=oneshot
ExecStart="/bin/bash /usr/local/bin/wp-scan.sh"
