# couchpotato/templates/couchpotato.service.erb
[Unit]
Description=CouchPotato application instance
After=network.target

[Service]
ExecStart=<%= scope['couchpotato::install_dir'] %>/CouchPotato.py --daemon
GuessMainPID=no
Type=forking

[Install]
WantedBy=multi-user.target