# couchpotato/templates/couchpotato.service.erb
[Unit]
Description=CouchPotato application instance
After=network.target

[Service]
ExecStart=<%= scope['couchpotato::install_dir'] %>/CouchPotato.py --daemon
GuessMainPID=no
Type=forking
User=<%= scope['couchpotato::settings']['user'] %>
Group=<%= scope['couchpotato::settings']['user'] %>

[Install]
WantedBy=multi-user.target
