# couchpotato/templates/couchpotato.service.erb
[Unit]
Description=CouchPotato application instance
After=network.target

[Service]
ExecStart=<%= @install_dir %>/CouchPotato.py --daemon
GuessMainPID=no
Type=forking
User=<%= @settings['user'] %>
Group=<%= @settings['user'] %>

[Install]
WantedBy=multi-user.target
