# zabbix/templates/zab-part.cron.erb
#!/bin/sh
/usr/bin/mysql -h <%= scope['zabbix::server::settings']['dBHost'] %> -u <%= scope['zabbix::server::settings']['dBUser'] %> -p<%= scope['zabbix::server::settings']['dBPassword'] %> <%= scope['zabbix::server::settings']['dBName'] %> -e "CALL partition_maintenance_all('zabbix');"
