# zabbix/templates/zab-part.cron.erb
#!/bin/bash
# global variables
CUR_DATE=$(/bin/date +%s)
NEXT_DST_DATE=$(/bin/date +%s -d $(/usr/bin/timedatectl status | /bin/grep Next -A 2 | /usr/bin/tail -1 | /usr/bin/awk '{print $2}'))
LAST_DST_DATE=$(/bin/date +%s -d $(/usr/bin/timedatectl status | /bin/grep Last -A 2 | /usr/bin/tail -1 | /usr/bin/awk '{print $2}'))

# Calculate days
let DAYS_BETW_DATES=($NEXT_DST_DATE-$CUR_DATE)/86400
let TOT_DAY_BET_DT=($NEXT_DST_DATE-$LAST_DST_DATE)/86400
let WEEK_LESS_DIFF_DT=($TOT_DAY_BET_DT-8)

function clean_up_partitions(){
	/usr/bin/mysql -h <%= scope['zabbix::server::settings']['dBHost'] %> -u <%= scope['zabbix::server::settings']['dBUser'] %> -p<%= scope['zabbix::server::settings']['dBPassword'] %> <%= scope['zabbix::server::settings']['dBName'] %> -e "CALL partition_maintenance_all('zabbix');"
}

if [ $DAYS_BETW_DATES -lt 8 -o $DAYS_BETW_DATES -gt $WEEK_LESS_DIFF_DT ]; then
    echo "Day Savings is changing in next couple weeks"
    echo "Not running partition clean up"
else
    echo "Day Savings is not changing in next couple of weeks";
    echo "Running partition clean up"
    clean_up_partitions
fi
