# transmission/templates/transmission-daemon.service.erb
[Unit]
Description=Transmission BitTorrent Daemon
After=network.target

[Service]
User=<%= @settings['user'] %>
Type=notify
ExecStart=/usr/bin/transmission-daemon -f -g <%= @config_dir %> -e <%= @log_dir %>/<%= @log_file %> --log-debug
ExecReload=/bin/kill -s HUP $MAINPID

[Install]
WantedBy=multi-user.target
