# transmission/templates/transmission-daemon.service.erb
[Unit]
Description=Transmission BitTorrent Daemon
After=network.target

[Service]
User=<%= scope['transmission::server::settings']['user'] %>
Type=notify
ExecStart=/usr/bin/transmission-daemon -f -g <%= scope['transmission::server::config_dir'] %> -e <%= scope['transmission::server::log_dir'] %>/<%= scope['transmission::server::log_file'] %> --log-debug
ExecReload=/bin/kill -s HUP $MAINPID

[Install]
WantedBy=multi-user.target
