# zabbix/files/zabbix-agent.xml
<?xml version="1.0"?> 
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1"> 
<service_bundle type='manifest' name='zabbix-agent'> 
  <service name='application/zabbix-agent' type='service' version='1'>
    <create_default_instance enabled='false' /> 
    <single_instance />
    <dependency name="fs-local" grouping="require_all" restart_on="none" type="service">
	<service_fmri value="svc:/system/filesystem/local"/>
    </dependency>
    <exec_method type='method' name='start' exec='/lib/svc/method/svc-zabbix-agent start' timeout_seconds='60'/>
    <exec_method type='method' name='stop' exec=':kill' timeout_seconds='60' />
    <exec_method type='method' name='restart' exec='/lib/svc/method/svc-zabbix-agent restart' timeout_seconds='60' /> 
    <stability value='Stable' /> 
    <template> 
      <common_name>
        <loctext xml:lang='C'> Zabbix-Agent 2.0.9 </loctext> 
      </common_name>
      <documentation> 
        <manpage title='zabbix-agentd' section='8' manpath='/usr/local/share/man' /> 
      </documentation>
    </template>
  </service>
</service_bundle>
